# recipes-core/rauc/files/rauc-mark-good.service
[Unit]
Description=Rauc Good-marking Service

[Service]
Type=oneshot
ExecStart=@BINDIR@/rauc status mark-good

[Install]
WantedBy=multi-user.target

